# web/jammon-web.conf
#!upstart

description "Upstart script for Jammon web daemon"

start on startup
stop on shutdown

script
    exec /srv/jammon/web/main.js
end script

pre-stop script
    sleep 5
end script
